# lokinet/contrib/lokinet-firewall.crontab
# lokinet firewall cronjobs

SHELL=/bin/bash

0 * * * * root /usr/local/bin/lokinet-update-firewall.sh

